# deploy/unicornx-backend.service
[Unit]
Description=UnicornX Aleo Backend (Express.js on port 5170)
Documentation=https://github.com/egorble/unicornxaleo
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=@USER@
Group=@USER@
WorkingDirectory=@HOME@/server
ExecStart=/usr/bin/node index.js
EnvironmentFile=@HOME@/.env
Environment=NODE_ENV=production
Environment=PORT=5170

# Restart policy
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Graceful shutdown
KillSignal=SIGTERM
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=@HOME@/server/data @HOME@/server/logs @HOME@/server/db @HOME@/.aleo @HOME@/contracts

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unicornx-backend

[Install]
WantedBy=multi-user.target
